# rootfs/etc/services.d/gmc320/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start gmc3xx service
# ==============================================================================
CONFIG_PATH=/data/options.json

PORT="$(bashio::config 'port')"
BAUDRATE="$(bashio::config 'baudrate')"
MQTT_SERVER="$(bashio::config 'mqtt_server')"
MQTT_PORT="$(bashio::config 'mqtt_port')"
MQTT_USER="$(bashio::config 'mqtt_user')"
MQTT_PASSWORD="$(bashio::config 'mqtt_password')"
DELAY="$(bashio::config 'repeat')"

bashio::log.info "Starting gmc3xx service"

# Get device information
MODEL=$(gmc320 ${PORT} ${BAUDRATE} | jq -r .model)
bashio::log.info "Device model: ${MODEL}"
SERIAL=$(gmc320 ${PORT} ${BAUDRATE} | jq -r .serial)
bashio::log.info "Device serial number: ${SERIAL}"

# Send homeassistant MQTT discovery message
bashio::log.info "Publishing MQTT config topics"
mosquitto_pub -u ${MQTT_USER} -P ${MQTT_PASSWORD} -p ${MQTT_PORT} -h ${MQTT_SERVER} -t "homeassistant/sensor/gmc3xx_${SERIAL}/" -m ''
mosquitto_pub -u ${MQTT_USER} -P ${MQTT_PASSWORD} -p ${MQTT_PORT} -h ${MQTT_SERVER} -t "homeassistant/sensor/gmc3xx_${SERIAL}/config" -m "{\"name\": \"${MODEL} ${SERIAL: -4}\", \"icon\": \"mdi:radioactive\", \"unit_of_measurement\": \"CPM\", \"state_topic\": \"homeassistant/sensor/gmc3xx_${SERIAL}/state\",  \"value_template\": \"{{ value_json.cpm }}\"}"

# Send the readings
bashio::log.info "Start sending data"
while true; do
	bashio::log.info "Fetching data from device"
	DATA=$(gmc320 ${PORT} ${BAUDRATE})

	bashio::log.info "Publishing state data"
	mosquitto_pub -u ${MQTT_USER} -P ${MQTT_PASSWORD} -p ${MQTT_PORT} -h ${MQTT_SERVER} -t "homeassistant/sensor/gmc3xx_${SERIAL}/state" -m "${DATA}"

	bashio::log.info "Sleeping"
	sleep ${DELAY}
done
